# packaging/dmj-vault-apps-admin/debian/postinst
#!/bin/bash
set -e

case "$1" in
    configure)
        mkdir -p /etc/dmj-vault-apps-admin/conf.d

        if [ ! -f /etc/dmj-vault-apps-admin/conf ]; then
            cat > /etc/dmj-vault-apps-admin/conf <<'EOF'
# DMJ Vault Admin App Configuration
PORT=9701
HOST=127.0.0.1
WORKERS=2
EOF
        fi

        systemctl daemon-reload
        systemctl enable dmj-vault-apps-admin.service
        systemctl start dmj-vault-apps-admin.service || true
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ;;
esac

#DEBHELPER#
